# spec/fixtures/riak.plist
<?xml version="1.0" encoding="UTF-8"?>

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>dev.riak</string>

    <key>Boxen</key>
    <dict>
      <key>Description</key>
      <string>A key-value store with delusions of grandeur.</string>

      <key>Port</key>
      <integer>18098</integer>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/boxen/homebrew/bin/riak</string>
      <string>start</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <false />

    <key>UserName</key>
    <string>wfarr</string>

    <key>WorkingDirectory</key>
    <string>/opt/boxen/data/riak</string>

    <key>StandardErrorPath</key>
    <string>/opt/boxen/log/riak/console.log</string>

    <key>StandardOutPath</key>
    <string>/opt/boxen/log/riak/console.log</string>

    <key>HardResourceLimits</key>
    <dict>
      <key>NumberOfFiles</key>
      <integer>4096</integer>
    </dict>
  </dict>
</plist>
